# config/supervisor-sample.conf
[program:remote-dl-py]
directory=/sites/remote-dl-py
command=python app.py
autorestart=true
user=nginx
